# systemd/user/read-clock.timer
[Unit]
Description=Read current hour every hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
